# cleanup-mastodon.service
[Unit]
Description=Cleanup Mastodon
[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/cleanup-mastodon.sh
User=mastodon