# cleanup-mastodon.timer
[Unit]
Description=Clean up Mastodon weekly

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target